# parser/cluster/history-server.conf
[supervisord]
nodaemon=true

[program:spark-master]
command=/opt/spark/sbin/start-history-server.sh