# parser/cluster/slave.conf
[supervisord]
nodaemon=true

[program:spark-slave]
command=/opt/spark/sbin/start-slave.sh spark://spark-master:7077